# tvheadend/rootfs/etc/s6-overlay/s6-rc.d/ssh-server/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: TVHeadend
# Starts TVHeadend
# ==============================================================================

bashio::log.info "Starting the SSH daemon..."

mkdir -p /root/.ssh/
echo "" > /root/.ssh/authorized_keys

if bashio::config.has_value 'allowed_ssh_keys'; then
    while read -r arg; do
        echo "${arg}" >> /root/.ssh/authorized_keys
    done <<< "$(bashio::config 'allowed_ssh_keys')"
fi
cd /etc/ssh/
/usr/bin/ssh-keygen -A
exec /usr/sbin/sshd -D -e < /dev/null